# docker/crontab
* * * * * node /index.js
# crontab requires an empty line at the end of the file
